# ubuntu-net-rootfs/var/lib/dpkg/info/netscript-2.4.list
/.
/etc
/etc/init.d
/etc/init.d/netscript
/etc/insserv.conf.d
/etc/insserv.conf.d/netscript
/etc/netscript
/etc/netscript/README
/etc/netscript/if.conf
/etc/netscript/ipfilter.conf
/etc/netscript/network.conf
/etc/netscript/network.conf.d
/etc/netscript/network.conf.d/00default-eth0-dhcp.conf
/etc/netscript/qos.conf
/etc/netscript/wep.conf
/etc/ppp
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/netscript
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/netscript
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/netscript.service
/lib/systemd/system/netscript@.service
/lib/udev
/lib/udev/netscript-hotplug
/lib/udev/rules.d
/lib/udev/rules.d/85-netscript.rules
/sbin
/sbin/netscript
/usr
/usr/share
/usr/share/doc
/usr/share/doc/netscript-2.4
/usr/share/doc/netscript-2.4/2.4.x-kernel.config.gz
/usr/share/doc/netscript-2.4/BUGS
/usr/share/doc/netscript-2.4/NEWS.Debian.gz
/usr/share/doc/netscript-2.4/README
/usr/share/doc/netscript-2.4/README.Debian.gz
/usr/share/doc/netscript-2.4/README.etc
/usr/share/doc/netscript-2.4/TODO
/usr/share/doc/netscript-2.4/changelog.gz
/usr/share/doc/netscript-2.4/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/netscript-2.4
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/netscript-2.4.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/netscript.8.gz
/sbin/ifdown
/sbin/ifup
/usr/share/man/man5/if.conf.5.gz
/usr/share/man/man5/ipfilter.conf.5.gz
/usr/share/man/man5/network.conf.5.gz
/usr/share/man/man5/qos.conf.5.gz
/usr/share/man/man8/ifdown.8.gz
/usr/share/man/man8/ifup.8.gz
